# my_server/mysql/supervisor.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld
autostart=true 
autorestart=true